# process/cron/orders.cron
#!/bin/bash

## this script processes incoming order files.
# files are delivered into an incoming queue by procmail, then cron runs
# this here script to make a non-blocking syntax check and reject or
# accept the order file.

if [ "yes" != "$CONFIRM" ] ; then
  exit
fi

for GAME in $*
do
  $HOME/eressea/orders-php/check-orders.sh $GAME
done
