# process/cron/orders.cron
#!/bin/bash

## this script processes incoming order files.
# files are delivered into an incoming queue by procmail, then cron runs
# this here script to make a non-blocking syntax check and reject or
# accept the order file.

if [ "no" == "$CONFIRM" ] ; then
  exit
fi

if [ -z "$ERESSEA" ] ; then
  ERESSEA="$HOME/eressea"
  echo "The ERESSEA environment variable is not set. Assuming $ERESSEA."
fi

cd "$ERESSEA" || exit

for GAME in "$@"
do
  orders-php/check-orders.sh "$GAME"
done

