# process/cron/preview.cron
#!/bin/bash

[ "$PREVIEW" != "yes" ] && exit
[ -z ${ERESSEA} ] && ERESSEA=$HOME/eressea

branch="develop"
if [ -e ${ERESSEA}/build/.preview ]; then
    branch=`cat ${ERESSEA}/build/.preview`
fi
SRC=${ERESSEA}/git
${SRC}/s/preview build ${branch} || exit $?
for game in 2 3 4 ; do
	${SRC}/s/preview -g ${game} run && \
	${SRC}/s/preview -g ${game} send
done
